# supervisord.conf
[supervisord]
nodaemon=true

[program:freeswitch]
autorestart=true
command=/bin/sh /usr/bin/start-freeswitch.sh

[program:cpbx]
autorestart=true
directory=/opt
command=/opt/cpbx